# bksrv-message/src/bksrv-message.target
[Unit]
Description=bksrv-message target
Requires=bksrv-message@.service

[Install]
WantedBy=multi-user.target
