# bksrv-client/src/bksrv-client.target
[Unit]
Description=bksrv-client target.
Requires=bksrv-message.target

# zfs
Requires=zfs.target

[Install]
WantedBy=multi-user.target
