# bksrv-client/src/bksrv-client.timer
[Unit]
Description=bksrv-client timer

[Timer]
Unit=bksrv-client.service
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600
