# bksrv-send/src/bksrv-send.socket
[Unit]
Description=bksrv-send socket

[Socket]
Accept=yes
ListenStream=9999

[Install]
WantedBy=multi-user.target
