# bksrv-snapshot/src/bksrv-snapshot-hourly.timer
[Unit]
Description=bksrv-snapshot-hourly timer

[Timer]
Unit=bksrv-snapshot@hourly.service
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=59
