# bksrv-snapshot/src/bksrv-snapshot-daily.timer
[Unit]
Description=bksrv-snapshot-daily timer

[Timer]
Unit=bksrv-snapshot@daily.service
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600
